# data/systemd/ask-modules-load.service
[Unit]
Description=Load ASK fast-path kernel modules (cdx, auto_bridge, fci)
DefaultDependencies=no
Before=systemd-modules-load.service cmm.service
After=systemd-udevd.service
ConditionPathExists=/usr/local/lib/ask-modules/cdx.ko

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ask-modules-load.sh
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target